# conf/ghost.conf
[program:ghost]
command = /root/.nvm/v0.10.32/bin/node /var/www/ghost/index.js
directory = /var/www/ghost
user = root
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/ghost.log
stderr_logfile = /var/log/supervisor/ghost_err.log
environment = NODE_ENV="production"
